# rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
#umount /dev/mtdblock4 >& /dev/null
#umount /dev/overlayfs >& /dev/null
#umount /dev/sda1 >& /dev/null
#mount /dev/sda1 /overlay

#touch rc.sda1.worked
##
rm /etc/resolv.conf
cp /etc/resolv.conf.opendns /etc/resolv.conf
cp /root/resolv.conf.opendns /etc/resolv.conf

exit 0
