# root/etc-sample/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
/bin/opkg update && echo "## opkg cache refreshed ##" > /dev/tty0
exit 0
